# backup/backup-home.timer
[Unit]
Description=backup timer

[Timer]
OnCalendar=*-*-* 9:30:00

[Install]
WantedBy=timers.target
